# dist/upstart/cleversocks.conf
# CleverSocks SOCKS5 Proxy - Upstart job configuration
#
# Start:   sudo start cleversocks
# Stop:    sudo stop cleversocks
# Status:  sudo status cleversocks
# Logs:    /var/log/upstart/cleversocks.log

description "CleverSocks SOCKS5 Proxy"
author "clevotec"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

# Drop privileges. Create this user/group during package install:
#   useradd --system --no-create-home --shell /usr/sbin/nologin cleversocks
setuid cleversocks
setgid cleversocks

# File descriptor limit
limit nofile 65536 65536

exec /usr/bin/cleversocks -c /etc/cleversocks/config.yaml
